# prime-daemon/prime-daemon.service
[Unit]
Description=Prime Daemon - laptop control API for the Prime Android app
After=network-online.target tailscaled.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/Projects/prime/prime-daemon
ExecStart=%h/Projects/prime/prime-daemon/run.sh
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
